# conf/get_tor_nodes.service
[Unit]
Description=refresh tor node list

[Service]
Type=simple
ExecStart=/usr/local/bin/get_tor_nodes.sh 

