# conf/get_tor_nodes.service
[Unit]
Description=Backup

[Service]
Type=simple
ExecStart=/usr/local/bin/get_tor_nodes.sh 

